# daemonize.ini
[program:energy]
directory=/data/energy
command=/usr/bin/php serve.php
user=root
autostart=true
autorestart=true
startsecs=2
startretries=2
stdout_logfile_maxbytes=1MB
stderr_logfile=/data/energy/runtime/energy-err.log
stdout_logfile=/data/energy/runtime/energy-out.log
